# work/configuration/modules/kvm_2_0/templates/rc.local.erb
#!/bin/sh -e
for s in rx tx sgo tso ufo gso
do
  ethtool -K <%= cluster_iface -%> $s off 
done
exit 0
